# etc/profile.d/nss-sdb-cache.sh
# Fix NSS memory leak in SSL Curl requests
NSS_SDB_USE_CACHE=YES
export NSS_SDB_USE_CACHE
